# .config/autostart/mopidy1.desktop
[Desktop Entry]
Type=Application
Name=mopidy
Exec=mopidy
#X-GNOME-Autostart-Delay=5
StartupNotify=false
